# docker/novnc/conf.d/fluxbox.conf
[program:fluxbox]
command=/bin/bash -c '[[ "$MODE" == "primary" ]] && openbox'
autorestart=true
