# docker/novnc/conf.d/tint2.conf
[program:tint2]
command=/bin/bash -c '[[ "$MODE" == "primary" ]] && tint2'
autorestart=true
