# docker/novnc/conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc 3009 localhost:5900
autorestart=true
